# workshops/eda-workshop-soca/ansible/roles/create_users_groups/templates/etc/cron.d/create_users_groups
MAILTO=''
PATH="{{ScriptsDir}}:/sbin:/bin:/usr/sbin:/usr/bin"
*/5 * * * * root {{ScriptsDir}}/create_users_groups.py -i {{ConfigDir}}/users_groups.json
